# systemd/photovault-brightness.service
[Unit]
Description=PhotoVault Brightness Helper
After=network.target

[Service]
Type=simple
ExecStart=/home/ewastewa/photovault/scripts/brightness-helper.sh
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
